# backend/deploy/finance.service
[Unit]
Description=Finance Spring Boot API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/finance
EnvironmentFile=/etc/finance/finance.env
ExecStart=/usr/bin/java -jar /opt/finance/finance.jar
SuccessExitStatus=143
Restart=always
RestartSec=10

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/finance

[Install]
WantedBy=multi-user.target
